# systemd/chkbootal-stop.service
[Unit]
Description=Creates hashes

[Service]
Type=oneshot
ExecStart=/usr/bin/chkbootal.py shutdown
TimeoutSec=0
StandardOutput=tty
StandardInput=tty
RemainAfterExit=no

[Install]
WantedBy=shutdown.target
